# .init.script
#!/bin/sh

### BEGIN INIT INFO
<<<<<<< HEAD
# Provides:        brewerslabMaster
=======
# Provides:        piTempMonitorSlave
>>>>>>> 0b32c09d3a6959e0518743336a0f7c30e0e2a596
# Required-Start:  $network $remote_fs $syslog
# Required-Stop:   $network $remote_fs $syslog
# Default-Start:   2 3 4 5
# Default-Stop:
<<<<<<< HEAD
# Short-Description: Start piTempMonitor Master
=======
# Short-Description: Start piTempMonitor
>>>>>>> 0b32c09d3a6959e0518743336a0f7c30e0e2a596
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin

case $1 in
        start)
<<<<<<< HEAD
                sh /home/beer/brewerslab/master/start.sh
                status=$?
                ;;
        stop)
		sh /home/beer/brewerslab/master/stop.sh
=======
                sh /home/beer/brewerslab/slave/start.sh
                status=$?
                ;;
        stop)
		sh /home/beer/brewerslab/slave/stop.sh
>>>>>>> 0b32c09d3a6959e0518743336a0f7c30e0e2a596
                status=$?
                ;;
        restart)
                $0 stop && sleep 2 && $0 start
                ;;
        *)
                echo "Usage: $0 {start|stop|restart}"
                exit 2
                ;;
esac
